# 1.20-ubuntu-bionic/etc/supervisor/conf.d/register.conf
[program:register]
command=/usr/local/bin/register
stopsignal=KILL
redirect_stderr=true
autorestart=true
